# dist/etc/skel/.config/autostart/Picom Compositor.desktop
[Desktop Entry]
Categories=Utility;
Comment=An X compositor
Exec=picom
GenericName=X compositor
Icon=preferences-desktop-effects
Keywords=compositor;composite manager;window effects;transparency;opacity;
MimeType=
Name=Picom Compositor
NoDisplay=false
Path=
StartupNotify=true
Terminal=false
TerminalOptions=
TryExec=picom
Type=Application
Version=1.0
X-DBUS-ServiceName=
X-DBUS-StartupType=
X-KDE-SubstituteUID=false
X-KDE-Username=
